# deploy/gunicorn.conf
# Ubuntu upstart file at /etc/init/gunicorn.conf

description "Zheroes"
start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 10 5
exec /srv/zheroes/deploy/gunicorn.sh
